# config/50-z-udev.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/z-autoset-screen.py"
